# deploy/agentpay.service
[Unit]
Description=AgentPay — Payment Layer for AI Agents
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=/opt/agentpay
EnvironmentFile=/opt/agentpay/.env
ExecStart=/opt/agentpay/venv/bin/python server.py

Restart=always
RestartSec=5

# Безопасность
NoNewPrivileges=true
ProtectHome=false
ProtectSystem=full

# Лимиты
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
