# deploy/agentpay.service
[Unit]
Description=AgentPay — Payment layer for AI agents
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/agentpay
EnvironmentFile=/opt/agentpay/.env
ExecStart=/usr/bin/python3 server.py
Restart=always
RestartSec=5

# Безопасность
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=multi-user.target
